# config/supervisor/video_telegram.conf
[program:motion_alarm_telegram]
command=php /home/butyh/project/dvr/telegram_bot.php
directory=/home/butyh/project/dvr/
autorestart=true
redirect_stderr=true
loglevel=warn