# config/supervisor/video_cleanup.conf
[program:cleanup]
command=/home/butyh/project/dvr/cleanup_video.sh
directory=/home/butyh/project/dvr/
autorestart=true
redirect_stderr=true
loglevel=warn