# config/supervisor/video_rename.conf
[program:rename]
command=/home/pi/project/dvr/rename_video.sh
directory=/home/pi/project/dvr/
autorestart=true
redirect_stderr=true
loglevel=warn