# config/supervisor/video_capture.conf
[program:capture]
command=/home/butyh/project/dvr/capture_video.sh 192.168.0.33
directory=/home/butyh/project/dvr/
autorestart=true
redirect_stderr=true
loglevel=warn