# config/supervisor/video_cleanup.conf
[program:cleanup]
command=/home/pi/project/dvr/cleanup_video.sh
directory=/home/pi/project/dvr/
autorestart=true
redirect_stderr=true
loglevel=warn