# config/supervisor/video_motion_detect.conf
[program:motion_detect]
command=/home/butyh/project/dvr/motion_detect.sh 192.168.0.31 /home/butyh/frames
directory=/home/butyh/project/dvr/
autorestart=true
redirect_stderr=true
loglevel=warn